# wled-controller.service
[Unit]
Description=Norns WLED LED1/LED2 controller (demo_full.py)
After=network.target sound.target

[Service]
Type=simple
WorkingDirectory=%h/wled-controller
ExecStart=%h/wled-controller/venv/bin/python3 -u demo_full.py
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
